# usr/lib/ipkg/info/busybox.list
/bin/usleep
/usr/bin/tr
/usr/bin/clear
/usr/bin/hexdump
/etc/busybox.links
/usr/bin/wc
/bin/tar
/usr/bin/mkfifo
/bin/mknod
/usr/bin/dumpleases
/etc/udhcpc.d/50default
/usr/bin/find
/bin/uname
/usr/bin/cut
/usr/bin/fuser
/usr/bin/[
/bin/zcat
/usr/bin/bzcat
/usr/bin/head
/usr/bin/openvt
/sbin/insmod
/usr/sbin/udhcpd
/usr/bin/tee
/sbin/route
/usr/bin/loadfont
/usr/bin/strings
/sbin/setconsole
/bin/true
/bin/chmod
/usr/bin/printf
/etc/init.d/syslog
/bin/echo
/usr/bin/top
/bin/chgrp
/bin/vi
/usr/bin/dc
/usr/bin/du
/bin/cpio
/sbin/e2fsck
/usr/bin/diff
/sbin/udhcpc
/usr/bin/patch
/bin/rm
/usr/bin/ar
/bin/sh
/usr/bin/env
/bin/ps
/usr/bin/renice
/usr/sbin/rdate
/usr/bin/whoami
/usr/bin/time
/usr/bin/expr
/usr/bin/basename
/usr/bin/nc
/usr/bin/tty
/usr/bin/od
/usr/bin/test
/bin/fgrep
/bin/sleep
/usr/bin/id
/usr/bin/seq
/etc/syslog.conf
/usr/bin/reset
/usr/bin/uniq
/usr/bin/chvt
/sbin/mkfs.minix
/sbin/losetup
/bin/dd
/bin/df
/usr/bin/less
/bin/grep
/bin/cp
/usr/bin/nslookup
/usr/bin/killall
/bin/touch
/usr/bin/nohup
/bin/watch
/bin/mktemp
/sbin/switch_root
/usr/bin/yes
/sbin/swapoff
/bin/pidof
/bin/busybox
/usr/bin/unzip
/bin/cat
/sbin/ifup
/bin/mv
/sbin/ifconfig
/bin/egrep
/bin/ln
/bin/ls
/bin/date
/usr/bin/[[
/bin/ip
/bin/ping
/bin/umount
/sbin/syslogd
/sbin/mkswap
/usr/bin/which
/usr/bin/sort
/bin/sed
/sbin/mkfs.ext3
/sbin/mkfs.ext2
/sbin/ifdown
/bin/dumpkmap
/usr/bin/traceroute
/sbin/klogd
/sbin/loadkmap
/sbin/start-stop-daemon
/usr/bin/readlink
/bin/ash
/sbin/sysctl
/bin/mount
/bin/stty
/bin/gunzip
/usr/bin/uptime
/usr/share/udhcpc/default.script
/bin/kill
/sbin/pivot_root
/bin/chown
/usr/bin/awk
/usr/bin/who
/bin/chattr
/usr/bin/logname
/usr/bin/bunzip2
/bin/dmesg
/usr/bin/wget
/sbin/hwclock
/bin/sync
/sbin/swapon
/usr/bin/md5sum
/sbin/modprobe
/sbin/fdisk
/bin/gzip
/sbin/rmmod
/bin/ping6
/bin/hostname
/usr/bin/deallocvt
/sbin/fsck.minix
/usr/sbin/fbset
/bin/more
/sbin/mke2fs
/sbin/fsck.ext3
/sbin/lsmod
/sbin/fsck.ext2
/bin/false
/usr/bin/tail
/bin/pwd
/sbin/logread
/usr/bin/free
/usr/bin/realpath
/bin/rmdir
/etc/init.d/hwclock.sh
/bin/mkdir
/sbin/fsck
/bin/run-parts
/usr/sbin/chroot
/bin/netstat
/usr/bin/dirname
/usr/bin/xargs
