# usr/lib/ipkg/info/cron.postrm
	update-rc.d cron remove
update-rc.d $D cron remove
