# usr/lib/ipkg/info/cron.list
/usr/sbin/cron
/usr/bin/crontab
/etc/init.d/cron
