# usr/lib/ipkg/info/ipkg.postrm
#!/bin/sh
update-alternatives --remove ipkg /usr/bin/ipkg-cl
update-rc.d $D ipkg-configure remove
