# usr/lib/ipkg/info/lighttpd.list
/usr/bin/spawn-fcgi
/etc/init.d/lighttpd
/usr/sbin/lighttpd-angel
/etc/lighttpd.conf
/usr/sbin/lighttpd
