# usr/lib/ipkg/info/ntpclient.postrm
#!/bin/sh
update-rc.d $D ntpclient remove
