# usr/lib/ipkg/info/ntpclient.list
/etc/init.d/ntpclient
/sbin/ntpclient
/sbin/adjtimex
/usr/sbin/ntpclient-drift-rate.awk
